# leap-linux/board/raspberry_pi_4/rootfs_overlay/etc/systemd/system/leap.service
[Unit]
Description=leap

[Service]
ExecStart=/usr/bin/leap-server

[Install]
WantedBy=multi-user.target
