# leap-linux/board/raspberry_pi_4/rootfs_overlay/etc/systemd/system/leap.service
[Unit]
Description=main leap service
Wants=var-lib-leap.mount

[Service]
ExecStart=/usr/bin/leap-server --fallback
Type=exec
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
