# leap-linux/board/raspberry_pi_4/rootfs_overlay/etc/systemd/system/var-lib-leap.mount
[Unit]
Description=LEAP storage mount

[Mount]
What=LABEL=LEAP_DATA
Where=/var/lib/leap
Type=ext4
Options=rw,noatime
TimeoutSec=90
